# root/etc/systemd/system/gitea.service.d/nethserver-gitea.conf
[Service]
WorkingDirectory=/var/lib/nethserver/gitea
